# files/etc/supervisor.d/subversion-edge.conf
[program:csvn]
command=/opt/csvn/bin/csvn console 
user=collabnet
priority=1

[program:csvn-httpd]
command=/opt/csvn/bin/csvn-httpd start
priority=2
startretries=10